# systemd/nsmtp.socket
[Unit]
Description=Netty based SMTP server

[Socket]
ListenStream=25
Accept=no

[Install]
