# systemd/nsmtp.service
[Unit]
Description=Netty based SMTP server

[Service]
ExecStart=/usr/bin/java -jar /opt/nsmtp/nsmtp.jar
StandardInput=socket
